# scripts/com.davinomjr.investments-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent that runs scripts/push-portfolio.sh daily at 09:00.

  Install on a new machine:
    1. Edit the absolute path inside <ProgramArguments> below to match the
       location of this repo on the target Mac.
    2. cp scripts/com.davinomjr.investments-sync.plist \
         ~/Library/LaunchAgents/com.davinomjr.investments-sync.plist
    3. launchctl bootstrap gui/$(id -u) \
         ~/Library/LaunchAgents/com.davinomjr.investments-sync.plist
    4. (optional, needs sudo) wake the Mac shortly before 09:00:
         sudo pmset repeat wakeorpoweron MTWRFSU 08:55:00

  Run on demand (for testing):
    launchctl kickstart -k gui/$(id -u)/com.davinomjr.investments-sync

  This MUST be a LaunchAgent (in ~/Library/LaunchAgents) — not a Daemon.
  The B3 portal sits behind Cloudflare and rejects headless browsers, so the
  worker pops a real Chromium window, which requires a logged-in GUI session.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.davinomjr.investments-sync</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>/Users/davinomjr/dev/investments-portfolio-manager/scripts/push-portfolio.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/davinomjr/Library/Logs/investments-sync.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/davinomjr/Library/Logs/investments-sync.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>ProcessType</key>
    <string>Interactive</string>
</dict>
</plist>
